# browser/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:xvfb]
command=Xvfb :1 -screen 0 %(ENV_RESOLUTION)s -ac +extension RANDR
autorestart=true
priority=10

[program:x11vnc]
command=x11vnc -display :1 -nopw -forever -shared -listen 127.0.0.1 -rfbport 5900
autorestart=true
priority=20

[program:websockify]
command=websockify --web /usr/share/novnc 8080 127.0.0.1:5900
autorestart=true
priority=30

[program:chromium]
command=%(ENV_CHROMIUM_PATH)s --no-sandbox --remote-debugging-port=9223 --remote-debugging-address=127.0.0.1 --remote-allow-origins=* --user-data-dir=/tmp/chrome-user-data --window-size=1280,720 --force-device-scale-factor=1 --autoplay-policy=no-user-gesture-required --no-first-run --no-default-browser-check --disable-gpu --display=:1
autorestart=true
priority=40

[program:tcp-forwarder]
command=python3 /tcp_forward.py
autorestart=true
priority=50
